# wsprdaemon_server@.service
[Unit]
Description=WSPRDAEMON Server (%i)
After=network.target clickhouse-server.service
Wants=clickhouse-server.service

[Service]
Type=simple
User=wsprdaemon
Group=wsprdaemon

# Create required directories
ExecStartPre=/bin/mkdir -p /var/log/wsprdaemon
ExecStartPre=/bin/mkdir -p /var/lib/wsprdaemon
ExecStartPre=/bin/mkdir -p /tmp/wsprdaemon
ExecStartPre=/bin/chown -R wsprdaemon:wsprdaemon /var/log/wsprdaemon
ExecStartPre=/bin/chown -R wsprdaemon:wsprdaemon /var/lib/wsprdaemon
ExecStartPre=/bin/chown wsprdaemon:wsprdaemon /tmp/wsprdaemon
ExecStart=/usr/local/bin/wsprdaemon_server.sh /etc/wsprdaemon/%i.conf

Restart=on-failure
RestartSec=60

# Memory limits
MemoryMax=2G
MemoryHigh=1.5G

# Security hardening
NoNewPrivileges=true
PrivateTmp=false
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/wsprdaemon /var/lib/wsprdaemon /tmp/wsprdaemon /var/spool/wsprdaemon

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wsprdaemon-server-%i

[Install]
WantedBy=multi-user.target
