# wsprdaemon_reflector@.service
[Unit]
Description=WSPRDAEMON Reflector (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/home/wsprdaemon/wsprdaemon

# The %i is replaced by the instance name (e.g., "reflector")
# So this becomes: /usr/local/bin/wsprdaemon_reflector.sh /etc/wsprdaemon/reflector.conf
ExecStart=/usr/local/bin/wsprdaemon_reflector.sh /etc/wsprdaemon/%i.conf

# Restart policy
Restart=always
RestartSec=10

# Security hardening
NoNewPrivileges=true
PrivateTmp=false
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/lib/wsprdaemon
ReadWritePaths=/var/log/wsprdaemon
ReadWritePaths=/var/spool/wsprdaemon
ReadWritePaths=/home

# Resource limits
MemoryMax=512M
CPUQuota=50%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wsprdaemon-reflector-%i

[Install]
WantedBy=multi-user.target
