# wsprnet_scraper@.service
[Unit]
Description=WSPRNET Scraper Service (%i)
After=network.target clickhouse-server.service
Wants=clickhouse-server.service

[Service]
Type=simple
User=wsprdaemon
Group=wsprdaemon
ExecStart=/usr/local/bin/wsprnet_scraper.sh
Restart=on-failure
RestartSec=60

# Security hardening (relaxed for venv access)
NoNewPrivileges=true
PrivateTmp=true
ReadWritePaths=/var/log/wsprdaemon /var/lib/wsprdaemon

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wsprnet-scraper-%i

[Install]
WantedBy=multi-user.target
